# fogw.service
[Unit]
Description=Failover Gateway

[Service]
Type=oneshot
ExecStart=/opt/fogw/fogw.sh