# fogw.service
[Unit]
Description=fogw

[Service]
Type=oneshot
ExecStart=/opt/fogw/fogw.sh
